# Projects/Timbergrove/Petro-link/timbergrove_petrolink-tree_cep_engine_hg-700cbb5333f4/Engine/src/main/bin/mbe-linux-systemd.service
[Unit]
Description=Petrolink MBE Engine Service

[Service]
Type=idle
User=root

# The configuration file application.properties should be here:
WorkingDirectory=/opt/mbe
ExecStart=/usr/bin/java -Duser.dir=/opt/mbe -DlogPath=log -Dmq.data=amq -DIGNITE_PERFORMANCE_SUGGESTIONS_DISABLED=true -Xmx4G -XX:+HeapDumpOnOutOfMemoryError -cp /opt/mbe/petrolink-mbe-engine.jar:. com.petrolink.mbe.engine.ComplexEventProcessor
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target

